# infra/backup/crontab
# Sindireceita backup pipeline crontab — driven by supercronic inside the
# `backup` compose sidecar (see deploy/compose/compose.stg.yml).
#
# Schedule: 03:15 in the container's TZ (set via `environment.TZ` on the
# service; defaults to America/Sao_Paulo). Matches the legacy systemd-unit
# OnCalendar=*-*-* 03:15:00 America/Sao_Paulo timing.
#
# Each line: MIN HOUR DOM MONTH DOW  COMMAND
# Output (stdout + stderr) is captured by docker and forwarded to
# promtail/Loki. backup.sh emits its own structured key=value records on
# stderr (see ADR 0102 § "Consequences").
15 3 * * * /usr/local/bin/backup.sh
